# supervisord.conf
[supervisord]

nodaemon=true

[program:logstash]
command=/opt/start-logstash.sh

[program:proxy]
command=/opt/start-proxy.sh

[program:rss]
command=/opt/start-rss.sh